# docker/parrot-core/supervisor/supervisor-novnc.ini
[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --listen "%(ENV_NOVNC_PORT)s" --vnc "localhost:%(ENV_VNC_PORT)s"
redirect_stderr=true
stdout_logfile=/var/log/novnc/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/novnc/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
