# docker/parrot-core/supervisor/supervisor-ttyd.ini
[program:ttyd]
command=/usr/bin/ttyd --port "%(ENV_TTY_PORT)s" --credential "%(ENV_USERNAME)s:%(ENV_PASSWORD)s" /bin/bash
redirect_stderr=true
stdout_logfile=/var/log/ttyd/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/ttyd/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
